# files/auth/supervisord.conf
{{- $programs := "beat,discord,gunicorn" }}
[supervisord]
nodaemon=true
user=allianceauth

[program:beat]
command=/opt/venv/bin/celery -A myauth beat
directory=/home/allianceauth/myauth
user=allianceauth
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startsecs=10
priority=998
stdout_events_enabled=true
stderr_events_enabled=true

{{ if eq (int .Values.celery.workers) 1 }}
[program:worker]
command=/opt/venv/bin/celery -A myauth worker -l INFO --pool threads --concurrency {{ .Values.celery.threads }}
directory=/home/allianceauth/myauth
user=allianceauth
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=998
stdout_events_enabled=true
stderr_events_enabled=true
{{- $programs = printf "%s,%s" $programs "worker" }}
{{- end }}

[program:discord]
command=nice /opt/venv/bin/celery -A myauth worker --max-tasks-per-child=250 --concurrency=1 -Q discordqueue
directory=/home/allianceauth/myauth
user=allianceauth
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=998

[program:gunicorn]
user=allianceauth
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stderr
directory=/home/allianceauth/myauth
command=python /home/allianceauth/myauth/manage.py runserver --noreload 0.0.0.0:8000
autostart=true
autorestart=true
stopsignal=INT
stdout_events_enabled=true
stderr_events_enabled=true

{{ if .Values.services.discordbot.enabled }}
[program:authbot]
command=python /home/allianceauth/myauth/bot_conf.py
directory=/home/allianceauth/myauth
user=allianceauth
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=600
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
{{- $programs = printf "%s,%s" $programs "authbot" }}
{{- end }}

{{- if .Values.services.discordproxy.enabled }}
[program:discordproxy]
command=/opt/venv/bin/discordproxyserver --token "%(ENV_DISCORD_PROXY_BOT_TOKEN)s"
directory=/home/allianceauth/myauth
user=allianceauth
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=120
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
{{- end }}

[group:myauth]
programs={{ $programs }}
priority=999

[supervisorctl]